# deploy/crontab
*/5 * * * * echo "Hello world" | logger -t opsworks-sample-php
*/5 * * * * echo "Goodbye world" | logger -t opsworks-sample-php
